# deploy/atm.service
[Unit]
Description=ATM Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/atm
ExecStart=/usr/bin/node src/bot.js
Restart=always
RestartSec=10
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
